# forge/libinput-rs.service
[Service]
Name = "libinput-rs"
Description = "Companion input preprocessor (evdev grab + uinput); optional beside libinput"
Type = "daemon"
ExecStart = "/usr/bin/libinput-rs"
User = "root"
Group = "input"

[Dependencies]
Requires = ["udev.service"]
After = ["udev.service"]
Before = ["display-manager.service"]

[Restart]
Condition = "always"
Delay = "2s"
